# .config/supervisor_app.conf
[program:uwsgi]
command=pipenv run uwsgi --ini /srv/project/.config/uwsgi_socket.ini

[program:nginx]
command=nginx